# quagga-isep/init-isep.service
[Service]
ExecStart=/usr/local/sbin/init-isep.sh
[Install]
WantedBy=default.target

